# etc/systemd/honcho-inspector-ui.service
[Unit]
Description=Honcho Inspector UI (Angular dev server)
Documentation=https://github.com/cloudbsdorg/honcho-inspector-ui
After=network-online.target honcho-inspector.service
Wants=network-online.target honcho-inspector.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/usr/local/share/honcho-inspector-ui

Environment=NG_CLI_ANALYTICS=false
Environment=NODE_ENV=production
# Pin PATH so the apt-installed Node wins over any SDKMAN or
# homebrew-installed Node in /usr/local/bin. /usr/bin comes first;
# the deb depends on the distro's nodejs package (>= 20) which
# satisfies Angular CLI's minimum version.
Environment=PATH=/usr/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/sbin:/bin
# Pin npm's cache and config dirs to paths the service user owns.
# The system's www-data has HOME=/var/www by default; systemd
# ignores Environment=HOME=... for system users (a known systemd
# quirk). We work around it by setting npm_config_cache and
# XDG_CACHE_HOME explicitly, both of which npm and Node respect
# unconditionally.
Environment=npm_config_cache=/var/lib/honcho-inspector-ui/.npm
Environment=XDG_CACHE_HOME=/var/lib/honcho-inspector-ui/.cache

# The Angular dev server IS the public listener. It binds on all
# interfaces (0.0.0.0:4200). The app is self-contained: an operator
# can hit it directly from a browser without any fronting reverse
# proxy. The package does not ship, configure, or depend on one.
#
# The backend (honcho-inspector.service) is bound to loopback only;
# browser->/api/* requests go through the Vite proxy in proxy.conf.json
# to the backend's loopback listener.
#
# ExecStartPre installs node_modules via `npm ci`. This is a one-shot
# on first boot: npm is idempotent so subsequent boots just check
# the lock and skip. The dep cache lives under the working directory's
# .angular/cache and node_modules -- which are writable thanks to
# ReadWritePaths below. The first boot can take 30-60s for npm ci
# on a fresh host (downloading the dep tree + the arch-specific
# esbuild binary); TimeoutStartSec=300 gives it 5 minutes.
ExecStartPre=/usr/bin/env npm ci --omit=dev --no-audit --no-fund
ExecStart=/usr/bin/env npx --no-install ng serve \
    --host 0.0.0.0 \
    --port 4200 \
    --proxy-config proxy.conf.json \
    --watch=false
Restart=always
RestartSec=5
TimeoutStartSec=300
TimeoutStopSec=30

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
# ProtectProc=ptraceable (not invisible) -- ng serve needs to read
# /proc/net/dev to enumerate network interfaces for the dev server
# bind. `invisible` blocks /proc/<pid>/net even for the service
# itself, which breaks libuv's interface enumeration (errno 97).
ProtectProc=ptraceable
RestrictRealtime=true
RestrictSUIDSGID=true
LockPersonality=true
SystemCallArchitectures=native

# ReadWritePaths: the working dir needs .angular/cache/ (Vite dep
# optimization), node_modules/ (npm install target), and any esbuild
# temp files. ReadWritePaths must list all of them or ProtectSystem=strict
# blocks writes and the service dies with EACCES or ENOENT.
ReadWritePaths=/var/lib/honcho-inspector-ui /usr/local/share/honcho-inspector-ui

# RestrictAddressFamilies intentionally NOT set. ng serve's interface
# discovery via libuv's uv_interface_addresses calls getifaddrs(3),
# which iterates over AF_PACKET (link-layer) on Linux even when the
# caller only binds on AF_INET. AF_PACKET is not in systemd's
# whitelisted address families, so setting RestrictAddressFamilies
# at all causes ng serve to fail with errno 97 (EAFNOSUPPORT). The
# service is already bound to loopback semantics via the unit's
# `--host 0.0.0.0` flag and the nginx front-end (if any) -- the
# AF-family filter was belt-and-suspenders hardening we don't need.

CapabilityBoundingSet=
AmbientCapabilities=

StandardOutput=journal
StandardError=journal
SyslogIdentifier=honcho-inspector-ui

[Install]
WantedBy=multi-user.target